# scripts/o11y-server.service
[Unit]
Description=OpenClaw O11y Server
After=network.target

[Service]
User=root
# Change to the path where you cloned the repository
WorkingDirectory=/opt/clawo11y

# Start the Go server binary
ExecStart=/opt/clawo11y/bin/clawo11y-server

# Environment variables
Environment="O11Y_SERVER_ADDR=0.0.0.0:8000"
Environment="O11Y_DB_URL=sqlite:////opt/clawo11y/data/o11y_server.db"
# Environment="O11Y_DATA_RETENTION_DAYS=7"
# Environment="O11Y_SERVER_SHUTDOWN_TIMEOUT_SEC=5"

Restart=always
RestartSec=5
# Limit logs
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=clawo11y-server

[Install]
WantedBy=multi-user.target
